# card_reader_adaptor.conf
[program:device-controller]
user=root
command=/usr/bin/card_reader_adaptor_bin
autostart = true
autorestart = true

